# deploy/note-agent.timer
[Unit]
Description=note trend agent daily timer
Requires=note-agent@%i.service

[Timer]
# 毎日 08:00 JST (= UTC 23:00 前日) に実行
OnCalendar=*-*-* 23:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
